# skills/gateway-watchdog/scripts/com.openclaw.gateway-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.openclaw.gateway-watchdog</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/env</string>
		<string>python3</string>
		<string>OPENCLAW_WATCHDOG_SCRIPT</string>
	</array>
	<key>WorkingDirectory</key>
	<string>OPENCLAW_HOME</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>OPENCLAW_HOME</key>
		<string>OPENCLAW_HOME</string>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>OPENCLAW_HOME/logs/gateway-watchdog.out.log</string>
	<key>StandardErrorPath</key>
	<string>OPENCLAW_HOME/logs/gateway-watchdog.err.log</string>
</dict>
</plist>
